# deploy/daily-arxiv-scheduler.service
[Unit]
Description=Daily arXiv Scheduler Service / 调度器服务 (template)
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/daily-arxiv
Environment="PATH=/path/to/python/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
EnvironmentFile=-/path/to/daily-arxiv/.env
ExecStart=/path/to/python/bin/python /path/to/daily-arxiv/scheduler.py
Restart=always
RestartSec=30
StandardOutput=append:/path/to/daily-arxiv/logs/scheduler.log
StandardError=append:/path/to/daily-arxiv/logs/scheduler.error.log

[Install]
WantedBy=multi-user.target
